# utils/supervisor-app.conf
[program:julia-app]
command = /usr/local/bin/run.sh

[program:nginx-app]
command = /usr/sbin/nginx -c /home/DIVAnd/DIVAnd-REST/utils/nginx.conf